# scripts/service-templates/paperboy-watcher.service
# Paperboy Watcher — systemd user unit
#
# HOW TO USE:
# 1. Replace /path/to/node with your Node.js path (run: which node)
# 2. Replace /path/to/watch-entry.js with the full path to dist/watch-entry.js
# 3. Copy and enable:
#    cp paperboy-watcher.service ~/.config/systemd/user/
#    systemctl --user enable --now paperboy-watcher
#
# To check status:
#    systemctl --user status paperboy-watcher
#
# To view logs:
#    journalctl --user -u paperboy-watcher

[Unit]
Description=Paperboy Watcher — auto-send .md files to Kindle
After=network-online.target

[Service]
Type=simple
ExecStart=/path/to/node /path/to/watch-entry.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
